# supervisord-cassandra.ini
[program:contrail-database]
command=cassandra -f
