# supervisord-rabbitmq.ini
[program:rabbitmq]
command=/etc/init/rabbitmq.sh 
